# supervisord.conf
[supervisord]
nodaemon=true

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:ascloud]
command=/usr/local/bin/projector run android-studio
